# includes.container/driver/App/etc/systemd/system/xencelabs.service
[Unit]
Description=xencelabs multy user Service

[Service]
ExecStart=/usr/lib/xencelabs/restar.sh

[Install]
WantedBy=multi-user.target
